# systemd/clanker-linear-agent.service.template
[Unit]
Description=Clanker Linear Agent for pi
Documentation=file:/home/exedev/clanker-linear-agent/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/exedev/clanker-linear-agent
EnvironmentFile=/home/exedev/clanker-linear-agent/.env
Environment=TOKEN_STORE_PATH=/home/exedev/clanker-linear-agent/data/linear-tokens.json
Environment=STATE_STORE_PATH=/home/exedev/clanker-linear-agent/data/oauth-states.json
Environment=PATH=/home/exedev/.local/share/mise/installs/node/latest/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/exedev/.local/share/mise/installs/node/latest/bin/node dist/server.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
